# meta-renesas/meta-rcar-gen2/recipes-graphics/gles-module/gles-user-module/rc.pvr.service

[Unit]
SourcePath=/usr/local/bin/rc.pvr
Before=weston.service

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/local/bin/rc.pvr start
ExecStop=/usr/local/bin/rc.pvr stop

[Install]
RequiredBy=weston.service
